# scripts/local/com.cbusillo.every-code-worker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cbusillo.every-code-worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/code/scripts/local/every-code-worker.sh</string>
    <string>run</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LAUNCHPLANE_EVERY_CODE_SERVICE</key>
    <string></string>
    <key>LAUNCHPLANE_EVERY_CODE_WORKSPACE</key>
    <string>/ABSOLUTE/PATH/TO/workspace</string>
    <key>LAUNCHPLANE_EVERY_CODE_WORKTREE</key>
    <string>/ABSOLUTE/PATH/TO/workspace/launchplane</string>
    <key>LAUNCHPLANE_EVERY_CODE_STATE_DIR</key>
    <string>/ABSOLUTE/PATH/TO/state/launchplane/every-code</string>
    <key>LAUNCHPLANE_EVERY_CODE_TOKEN_ITEM</key>
    <string>launchplane-every-code-worker-token</string>
    <key>LAUNCHPLANE_EVERY_CODE_HOST</key>
    <string>LOCAL-HOST-NAME</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/code</string>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/state/launchplane/every-code/every-code-worker/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/state/launchplane/every-code/every-code-worker/launchd.err.log</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
